# services/ai-opponent/supervisord.conf
[supervisord]
nodaemon=true

[program:websocket]
command=/app/./myDQN

[program:http_factory]
command=/app/./myFactory 8080 9876